# X86_64-Lede/files/etc/uci-defaults/99-custom.sh
#!/bin/sh

# Set default theme to luci-theme-argon
uci set luci.main.mediaurlbase='/luci-static/design'
uci commit luci

exit 0
